# etc/supervisor/conf.d/slapd.conf
[program:slapd]
command=/usr/sbin/slapd -F /etc/ldap/slapd.d -d 0
